# supervisord.conf
[supervisord]
nodaemon=true

[program:demo-service]
command=/app/go-api